# ansible/roles/syncthing/templates/syncthing.service.j2
[Service]
# Allow binding to all interfaces
Environment="STHOMEDIR={{ syncthing_config_dir }}"
